# builds/sbo/gcc8/profile.d/gcc8.csh
setenv CC     gcc-8
setenv CPP    cpp-8
setenv CXX    g++-8
setenv AR     gcc-ar-8
setenv NM     gcc-nm-8
setenv RANLIB gcc-ranlib-8
