# builds/sbo/gcc8/profile.d/gcc8.sh
export CC=gcc-8
export CPP=cpp-8
export CXX=g++-8
export AR=gcc-ar-8
export NM=gcc-nm-8
export RANLIB=gcc-ranlib-8

